# etc/lupin/systemd/system/sockettest.socket
[Unit]
Conflicts=sockettest.service

[Socket]
ListenStream=1337
Accept=yes

[Install]
WantedBy=sockets.target
