# etc/lupin/systemd/system/getty@tty1.service.d/autologin.conf
[Unit]
Requires=dbus.service
After=dbus.service

[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin florian --noclear %I 38400 linux
Type=simple
