# etc/ginny/NetworkManager/dispatcher.d/ntp
#!/bin/bash
[[ "$2" == up ]] && /usr/bin/ntpd -q -g
