# etc/ginny/NetworkManager/dispatcher.d/ntp
#!/bin/bash
[[ "$2" == up ]] && sntp -s pool.ntp.org
