# etc/ginny/NetworkManager/dispatcher.d/sshtunnel
#!/bin/bash
[[ "$2" == up ]] && sshtunnel start-all
[[ "$2" == down ]] && sshtunnel stop-all
